# today.service
[Unit]
Description=Sets the wallpaper

[Service]
ExecStart= /home/tobias/bin/today

[Install]
WantedBy=window-manager.target
